# runit/proxmox-daily-update/run
#!/bin/bash

ln -sf /usr/lib/*/proxmox-backup/proxmox-daily-update /etc/cron.daily/

echo "CRON: Starting..."
exec cron -f
