# runit/proxmox-backup-api/run
#!/bin/bash

chown backup:backup /etc/proxmox-backup
chmod 700 /etc/proxmox-backup

echo "API: Starting..."
exec /usr/lib/aarch64-linux-gnu/proxmox-backup/proxmox-backup-api
